# systemd/.config/systemd/user/hyprpaper-shuffle.timer
[Unit]
Description=Shuffle wallpapers every 30 minutes

[Timer]
OnBootSec=5sec
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
